# debian/traced-perf.service
[Unit]
Description=Perfetto sampling and profiling data source

[Service]
ExecStart=/usr/sbin/traced_perf
User=root
Group=root
